# earthquakeapi/EarthquakeAPI.service
[Unit]
Description=EarthquakeAPI

[Service]
Type=notify
ExecStart=/usr/sbin/earthquakeapi

[Install]
WantedBy=multi-user.target